# groundwork-monitor-ganglia/scripts/services/check-ganglia/run
#!/bin/sh
sleep 10
exec 2>&1
exec setuidgid nagios /usr/local/groundwork/nagios/libexec/check_ganglia.pl
